# external/etc/systemd/system/focusrite-resume.service
[Unit]
Description=Reset Scarlett 18i20 after suspend
After=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix-focusrite.sh

[Install]
WantedBy=suspend.target
